# download-guardian/download-guardian.service
[Unit]
Description=Download Gaurdian (watches ~/Downloads and prompts to scan new files using ClamAV)
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/download-guardian-scan
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

